# examples/launchd/com.herald.brief.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.herald.brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/herald/modules/x-briefs/bin/run-brief.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/herald/modules/x-briefs</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 20:00 in the machine's local timezone. Adjust to taste. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>20</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/herald/modules/x-briefs/logs/launchd-brief.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/herald/modules/x-briefs/logs/launchd-brief.log</string>
</dict>
</plist>
