# examples/launchd/com.herald.listener.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.herald.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>core/listen_discord.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/herald</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/Users/YOUR_USERNAME/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- The one intentionally persistent process — a Discord gateway client,
         not a scheduled job. KeepAlive restarts it on crash or exit;
         ThrottleInterval stops a crash-loop from hammering Discord. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/path/to/herald/logs/launchd-listener.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/herald/logs/launchd-listener.log</string>
</dict>
</plist>
